# deployments/systemd/cloud-flow-agent.service
[Unit]
Description=Cloud Flow Agent - 云内流量监测探针
Documentation=https://github.com/meinanzilinzhengying/cloud-flow-agent
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
User=root
Group=root
ExecStart=/usr/local/bin/cloud-flow-agent -config /etc/cloud-flow/agent.yaml
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5
TimeoutStartSec=30
TimeoutStopSec=30

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096
LimitCORE=0

# 安全加固
NoNewPrivileges=false  # eBPF 需要 root 权限
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/var/log/cloud-flow /etc/cloud-flow

# 环境变量
Environment="GOMAXPROCS=1"
Environment="GOMEMLIMIT=512MiB"

# 标准输出
StandardOutput=journal
StandardError=journal
SyslogIdentifier=cloud-flow-agent

# 优雅关闭
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
